# cluster/cell/upstart/consul.conf
#!upstart

start on started bootstrap
stop on shutdown
respawn

pre-start script
    echo 'nameserver 127.0.0.1' > /etc/resolvconf/resolv.conf.d/head
    resolvconf -u
end script

script
    echo "UPSTART: Trying to start consul - `date --rfc-3339=ns`"
    export $(cat /var/lattice/setup/lattice-environment)

    consul agent \
        -config-file /var/lattice/config/consul.json \
        -data-dir /tmp/consul \
        -join $CONSUL_SERVER_IP \
        >> /var/lattice/log/consul-service.log 2>&1
end script

post-stop exec sleep 5
