# cluster/brain/upstart/haproxy-reload.conf
#!upstart

start on file FILE=/var/lattice/config/haproxy/haproxy.cfg
stop on shutdown

script
 kill -HUP $(initctl status haproxy | cut -d ' ' -f 4)
end script