# coordinator/upstart/etcd.conf
#!upstart

start on started bootstrap

stop on shutdown
respawn

script
    echo "UPSTART: Trying to start etcd - `date --rfc-3339=ns`"
    etcd
end script

post-stop exec sleep 5