# coordinator/upstart/gnatsd.conf
#!upstart

start on started bootstrap
stop on shutdown
respawn

script
    echo "UPSTART: Trying to start gnatsd - `date +%s%N`"

    gnatsd -c /var/lattice/config/gnatsd.conf
end script
