# coordinator/upstart/gorouter.conf
#!upstart

start on started gnatsd
stop on shutdown
respawn

script
    echo "UPSTART: Trying to start gorouter - `date +%s%N`"
    gorouter -c /var/lattice/config/gorouter.yml
end script
