# coordinator/upstart/lattice-coordinator-bootstrap.conf
#!upstart

start on (local-filesystems and net-device-up IFACE!=lo)

pre-start script
    echo "UPSTART: Trying to start lattice coordinator service - `date +%s%N`"
    /var/lattice/setup/infrastructure-specific-setup
end script

script
 echo "Lattice Coordinator is ready!"
end script
