# coordinator/upstart/route-emitter.conf
#!upstart

start on started gnatsd
stop on shutdown
respawn

script
echo "UPSTART: Trying to start route-emitter - `date +%s%N`"
    route-emitter -debugAddr=0.0.0.0:17009
end script
