# diego-cell/upstart/bootstrap.conf
#!upstart

start on (local-filesystems and net-device-up IFACE!=lo)

pre-start script
    echo "UPSTART: Trying to start lattice bootstrap service - `date --rfc-3339=ns`"
end script
