# diego-cell/upstart/converger.conf
#!upstart

start on started lattice-diego-cell-bootstrap
stop on shutdown
respawn

script
    echo "UPSTART: Trying to start converger - `date +%s%N`"
    converger
end script
