# diego-cell/upstart/lattice-diego-cell-bootstrap.conf
#!upstart

start on (local-filesystems and net-device-up IFACE!=lo)

script
 echo "Lattice Diego Cell is ready!"
end script
