# diego-cell/upstart/lattice-diego-cell-bootstrap.conf
#!upstart

start on (local-filesystems and net-device-up IFACE!=lo)

pre-start script
    echo "UPSTART: Trying to start lattice coordinator service - `date +%s%N`"
    /var/lattice/setup/diego-cell-bootstrap
end script

script
 echo "Lattice Diego Cell is ready!"
end script
