# diego-cell/upstart/rep.conf
#!upstart

start on started gnatsd
stop on shutdown
respawn

script
    echo "UPSTART: Trying to start rep - `date +%s%N`"

    LRPHOST=`hostname -I | awk '{ print $1 }'`
    rep -stack=lucid64 -executorURL=http://127.0.0.1:1700 -lrpHost=$LRPHOST -cellID=cell_z1-0
end script
