# upstart/auctioneer.conf
#!upstart

start on started gnatsd
stop on shutdown
respawn

script
    echo "UPSTART: Trying to start auctioneer - `date +%s%N`"
    auctioneer
end script
