# upstart/auctioneer.conf
#!upstart

start on filesystem
stop on shutdown
respawn

script
    auctioneer
end script
