# upstart/auctioneer.conf
#!upstart

start on filesystem
stop on shutdown

script
    auctioneer
end script
