# upstart/converger.conf
#!upstart

start on filesystem
stop on shutdown

script
    converger
end script
