# upstart/converger.conf
#!upstart

start on (local-filesystems and net-device-up IFACE!=lo)
stop on shutdown
respawn

script
    echo "UPSTART: Trying to start converger - `date +%s%N`"
    converger
end script
