# upstart/converger.conf
#!upstart

start on filesystem
stop on shutdown
respawn

script
    converger
end script
