# upstart/diego-edge-bootstrap.conf
#!upstart

start on (local-filesystems and net-device-up IFACE!=lo)

script
    echo "UPSTART: Trying to start diego-edge wrapper service - `date +%s%N`"
end script
