# upstart/diego-edge-bootstrap.conf
#!upstart

start on (local-filesystems and net-device-up IFACE!=lo)

export SYSTEM_DOMAIN

script
    echo "UPSTART: Trying to start diego-edge wrapper service - `date +%s%N`"
    echo "UPSTART: Detecting infrastructure..."
    aws_metadata_api_base_url=http://169.254.169.254/2014-02-25
    curl -s --connect-timeout 3 $aws_metadata_api_base_url
    aws_curl_result=$?
   if [ "$aws_curl_result" -eq 0 ]; then
     echo "Detected AWS"
     SYSTEM_DOMAIN=`curl --connect-timeout 3 $aws_metadata_api_base_url/meta-data/public-ipv4`.xip.io
     echo "Setting xip.io IP's for publicly accessible routes to $SERVER_IP retrieved from amazon"
   else
     echo "Detected Non-AWS"
     SYSTEM_DOMAIN=192.168.11.11.xip.io
     echo "Defaulting xip.io IP's for publicly accessible routes to 192.168.11.11"
   fi

   echo $SYSTEM_DOMAIN > /var/diego/system-domain
   sed -i "s/SYSTEMDOMAIN_PLACEHOLDER/$SYSTEM_DOMAIN/" /var/diego/config/trafficcontroller.json
end script
