# upstart/diego-edge-bootstrap.conf
#!upstart

start on (local-filesystems and net-device-up IFACE!=lo)

pre-start script
    echo "UPSTART: Trying to start diego-edge wrapper service - `date +%s%N`"
    /var/diego/setup/infrastructure-specific-setup
end script

script
 echo "Diego Edge is ready!"
end script
