# upstart/diego-edge-bootstrap.conf
#!upstart

start on (local-filesystems and net-device-up IFACE!=lo)

export SYSTEM_DOMAIN

pre-start script
    echo "UPSTART: Trying to start diego-edge wrapper service - `date +%s%N`"
    echo "UPSTART: Detecting infrastructure..."
    /var/diego/config/provision_diego

    SYSTEM_DOMAIN=`cat /var/diego/system-domain`
end script

script
 echo "Diego Edge is ready!"
end script
