# upstart/etcd.conf
#!upstart

start on filesystem
stop on shutdown
respawn

script
    etcd
end script
