# upstart/etcd.conf
#!upstart

start on started lattice-bootstrap

stop on shutdown
respawn

script
    echo "UPSTART: Trying to start etcd - `date +%s%N`"
    etcd
end script
