# upstart/file-server.conf
#!upstart

start on filesystem
stop on shutdown
respawn

script
    file-server -port=8070 -staticDirectory=/var/diego/circus -skipCertVerify=true -ccAddress=IGNORED -ccUsername=IGNORED -ccPassword=IGNORED
end script
