# upstart/gnatsd.conf
#!upstart

start on started diego-edge-bootstrap
stop on shutdown
respawn

script
    echo "UPSTART: Trying to start gnatsd - `date +%s%N`"

    gnatsd -c /var/diego/config/gnatsd.conf
end script
