# upstart/gnatsd.conf
#!upstart

start on filesystem
stop on shutdown
respawn

script
    gnatsd -c /var/diego/config/gnatsd.conf
end script
