# upstart/gorouter.conf
#!upstart

start on filesystem
stop on shutdown

script
    gorouter -c /var/diego/config/gorouter.yml
end script
