# upstart/lattice-bootstrap.conf
#!upstart

start on (local-filesystems and net-device-up IFACE!=lo)

pre-start script
    echo "UPSTART: Trying to start lattice wrapper service - `date +%s%N`"
    /var/lattice/setup/infrastructure-specific-setup
end script

script
 echo "Diego Edge is ready!"
end script
