# upstart/rep.conf
#!upstart

start on filesystem
stop on shutdown
respawn

script
    rep -stack=lucid64 -executorURL=http://127.0.0.1:1700 -lrpHost=10.0.2.15 -executorID=cell_z1-0
end script
