# upstart/route-emitter.conf
#!upstart

start on filesystem
stop on shutdown
respawn

script
    route-emitter -debugAddr=0.0.0.0:17009
end script
